# triboka-agro/triboka-flask.service
[Unit]
Description=Triboka Web3 Flask Application
After=network.target

[Service]
Type=simple
User=rootpanel
Group=rootpanel
WorkingDirectory=/home/rootpanel/web/app.triboka.com/backend
Environment=PATH=/usr/bin:/usr/local/bin
Environment=FLASK_APP=app_web3.py
Environment=FLASK_ENV=production
Environment=FLASK_DEBUG=False
Environment=PYTHONPATH=/home/rootpanel/web/app.triboka.com/backend
ExecStart=/home/rootpanel/web/app.triboka.com/start_triboka_web3.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target