# triboka-agro/triboka-erp-nextjs.service
[Unit]
Description=Triboka ERP Next.js Frontend
After=network.target triboka-erp-backend.service
Requires=triboka-erp-backend.service

[Service]
Type=simple
User=rootpanel
Group=rootpanel
WorkingDirectory=/home/rootpanel/web/app.triboka.com/triboka-erp/frontend
Environment=PATH=/usr/bin:/usr/local/bin:/usr/local/node/bin
Environment=NODE_ENV=production
Environment=PORT=3002
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target